# seafile/supervisord.conf
[supervisord]
nodaemon=true

[program:seafile]
command=/opt/seafile/seafile-start.sh

[program:seahub]
command=/opt/seafile/seahub-start.sh
startsecs = 10